# ntfy/rootfs/etc/s6-overlay/s6-rc.d/init-ntfy/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# init-ntfy (oneshot)
#
# 1. Renders /etc/ntfy/server.yml from the add-on options.
# 2. Optionally provisions an admin user in the persistent auth database.
#
# All persistent state (auth db, message cache, attachments) lives under /data,
# which Home Assistant keeps across restarts and updates.
# ==============================================================================

readonly CONFIG_FILE="/etc/ntfy/server.yml"
readonly DATA_DIR="/data"
readonly CACHE_FILE="${DATA_DIR}/cache.db"
readonly AUTH_FILE="${DATA_DIR}/user.db"
readonly ATTACHMENT_DIR="${DATA_DIR}/attachments"

declare base_url
declare upstream_base_url
declare auth_default_access
declare cache_duration
declare log_level
declare behind_proxy
declare auth_enabled

base_url=$(bashio::config 'base_url')
upstream_base_url=$(bashio::config 'upstream_base_url')
auth_default_access=$(bashio::config 'auth_default_access')
cache_duration=$(bashio::config 'cache_duration')
log_level=$(bashio::config 'log_level')

# -----------------------------------------------------------------------------
# Ensure persistent directories exist
# -----------------------------------------------------------------------------
mkdir -p /etc/ntfy "${ATTACHMENT_DIR}"

# -----------------------------------------------------------------------------
# Render /etc/ntfy/server.yml
# -----------------------------------------------------------------------------
bashio::log.info "Generating ntfy server configuration..."

{
    # base-url MUST match the canonical public URL that clients use to reach
    # this server (reverse proxy / Cloudflare Tunnel). It drives topic links,
    # attachment URLs and iOS/Firebase push. See DOCS.md for the ingress note.
    echo "base-url: \"${base_url%/}\""

    # ntfy listens on container port 80. Home Assistant ingress connects here,
    # and it is also published on the host (default 8199) for the reverse proxy.
    echo "listen-http: \":80\""

    # Trust X-Forwarded-* headers from the ingress proxy / reverse proxy.
    if bashio::config.true 'behind_proxy'; then
        echo "behind-proxy: true"
    else
        echo "behind-proxy: false"
    fi

    # Persistent message cache
    echo "cache-file: \"${CACHE_FILE}\""
    echo "cache-duration: \"${cache_duration}\""

    # Persistent attachment cache
    echo "attachment-cache-dir: \"${ATTACHMENT_DIR}\""

    # Authentication (persistent user database)
    if bashio::config.true 'auth_enabled'; then
        echo "auth-file: \"${AUTH_FILE}\""
        echo "auth-default-access: \"${auth_default_access}\""
        # Enable the web UI account page (sign in, manage access tokens). Without
        # this ntfy's web app shows no account menu — the "Users" under Settings
        # is only a client-side credential store, not server accounts. Signup
        # stays off (private server); users/tokens are minted by the admin.
        echo "enable-login: true"
    else
        bashio::log.warning "auth_enabled is false: server is OPEN (read-write for everyone)."
    fi

    # Upstream relay used to wake up iOS devices via ntfy.sh
    if bashio::var.has_value "${upstream_base_url}"; then
        echo "upstream-base-url: \"${upstream_base_url%/}\""
    fi

    # Logging
    echo "log-level: \"${log_level}\""
    echo "log-format: \"text\""

    # Keep websocket/JSON stream connections alive through proxies
    echo "keepalive-interval: \"45s\""
} > "${CONFIG_FILE}"

bashio::log.info "Wrote ${CONFIG_FILE}"

# -----------------------------------------------------------------------------
# Provision users in the persistent auth database (idempotent)
# -----------------------------------------------------------------------------
# All provisioning needs auth on and the auth db present. ntfy refuses to touch
# the db unless the file already exists; an empty file is enough for it to
# create the schema on first use.
if bashio::config.true 'auth_enabled'; then

    if [ ! -f "${AUTH_FILE}" ]; then
        touch "${AUTH_FILE}"
    fi

    # --- Admin user (role=admin, read-write to everything) ---
    if bashio::config.has_value 'admin_user' \
        && bashio::config.has_value 'admin_password'; then

        admin_user=$(bashio::config 'admin_user')
        admin_password=$(bashio::config 'admin_password')

        if ntfy user --config "${CONFIG_FILE}" list 2>/dev/null \
            | grep -qE "^user ${admin_user} "; then
            bashio::log.info "Admin user '${admin_user}' already exists; leaving it unchanged."
        else
            bashio::log.info "Creating admin user '${admin_user}'..."
            if NTFY_PASSWORD="${admin_password}" \
                ntfy user --config "${CONFIG_FILE}" add --role=admin "${admin_user}"; then
                bashio::log.info "Admin user '${admin_user}' created."
            else
                bashio::log.warning "Could not create admin user '${admin_user}'; continuing anyway."
            fi
        fi
    else
        bashio::log.info "No admin_user/admin_password set; skipping admin provisioning."
    fi

    # --- Extra scoped users (regular role) with per-user topic ACLs ---
    # Each entry provisions a login that only needs specific topics — e.g. a
    # phone that READS its notification topics (`hab-*`). The user is created
    # once; the password here is used only at creation (change it later with
    # `ntfy user change-pass`). The ACL is re-applied every start, so a topic
    # or access change in the options takes effect on the next restart.
    for idx in $(bashio::config 'users|keys'); do
        u_name=$(bashio::config "users[${idx}].name")
        u_password=$(bashio::config "users[${idx}].password")
        u_topics=$(bashio::config "users[${idx}].topics")
        u_access=$(bashio::config "users[${idx}].access")
        # bashio returns "null" for absent optional keys — normalise to empty.
        [ "${u_topics}" = "null" ] && u_topics=""
        [ "${u_access}" = "null" ] && u_access=""
        [ -n "${u_access}" ] || u_access="read-only"

        if ntfy user --config "${CONFIG_FILE}" list 2>/dev/null \
            | grep -qE "^user ${u_name} "; then
            bashio::log.info "User '${u_name}' already exists; leaving password unchanged."
        else
            bashio::log.info "Creating user '${u_name}'..."
            if ! NTFY_PASSWORD="${u_password}" \
                ntfy user --config "${CONFIG_FILE}" add "${u_name}"; then
                bashio::log.warning "Could not create user '${u_name}'; skipping."
                continue
            fi
        fi

        if [ -n "${u_topics}" ]; then
            bashio::log.info "Granting '${u_name}' ${u_access} on '${u_topics}'."
            ntfy access --config "${CONFIG_FILE}" "${u_name}" "${u_topics}" "${u_access}" \
                || bashio::log.warning "Could not set ACL for '${u_name}'."
        fi
    done
fi
